# configs/hyprland-dotfiles/airootfs/etc/systemd/system/plymouth-poweroff.service.d/display.conf
[Unit]
After=
After=plymouth-start.service
Conflicts=display-manager.service
