# configs/hyprland-dotfiles/airootfs/etc/systemd/system/plymouth-quit.service.d/wait-for-dm.conf
[Unit]
After=display-manager.service

[Service]
ExecStart=
ExecStart=-/usr/bin/plymouth quit --retain-splash
