# configs/hyprland-dotfiles/airootfs/etc/systemd/system/plymouth-quit.service.d/wait-for-dm.conf
[Service]
ExecStart=
ExecStartPre=/usr/bin/sleep 10
ExecStart=-/usr/bin/plymouth quit --retain-splash
