# configs/hyprland-dotfiles/airootfs/etc/systemd/system/plymouth-reboot.service.d/display.conf
[Unit]
After=
After=plymouth-start.service
Conflicts=display-manager.service

[Service]
ExecStartPost=/usr/bin/sleep 5
